# deploy/templates/calendar_sync.timer
[Unit]
Description=Timer for calendar_sync agent (every 15 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
